# docker/root/etc/s6-overlay/s6-rc.d/init-stravastat-config/up
/etc/s6-overlay/s6-rc.d/init-stravastat-config/run